# debian/bhp.cron.d
*/1 *	* * *	root	[ -x /usr/lib/bhp/cron ] && /usr/lib/bhp/cron
